# test-mta/supervisord.conf
[supervisord]
nodaemon=true

[program:fastapi]
command=uvicorn mail_store:app --host 0.0.0.0 --port %(ENV_MAIL_API_PORT)s
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:smtp_sink]
command=python smtp_sink.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=MAIL_API_PORT=%(ENV_MAIL_API_PORT)s,SMTP_PORT=%(ENV_SMTP_PORT)s

